# tools/battery-logger.service
[Unit]
Description=Battery Usage Logger
Documentation=https://github.com/travelingbear/power-profile-manager
After=multi-user.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /usr/local/bin/battery-logger.py
Restart=on-failure
RestartSec=60
StandardOutput=journal
StandardError=journal

# Run as user (not root)
User=%i
Group=%i

[Install]
WantedBy=default.target
