# src/power-profiled.service
[Unit]
Description=Power Profile Manager Daemon
Requires=tlp.service
After=multi-user.target tlp.service

[Service]
Type=simple
ExecStart=/usr/local/bin/power-profiled
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
